# templates/upstart.conf
description "genieacs-gui"
author "shweelan <shweelan@ubuntu-512mb-fra1-01>"

start on runlevel [2345]
stop on runlevel [016]

respawn

console log
script
    echo Fake service; sleeping for an hour...
    sleep 360
end script
